# 60-usb-libvirt.rules
ACTION=="bind|unbind", ENV{DEVTYPE}=="usb_device", SUBSYSTEM=="usb", RUN+="/etc/libvirt/libvirt-usb.sh"
