# 60-usb-libvirt.rules
ACTION=="bind" , SUBSYSTEM=="usb", RUN+="/etc/libvirt/libvirt-usb.sh"
